# systemd/tk-ssh-agent.service
[Unit]
Description=Trusted Key SSH Agent
RefuseManualStart=true

[Service]
ExecStart=/usr/bin/tk-ssh-agent agent --systemd
Type=simple
